# asset/supervisord.conf
[supervisord]
nodaemon = true
user = root

[program:mid]
command = /opt/agent/bin/mid.sh start
user = root